# debian/loom-tm.postinst
#!/bin/sh
systemctl enable loom-tm.service
